# deskflow-watchdog.service
[Unit]
Description=Deskflow Connection Watchdog
After=graphical-session.target

[Service]
Type=simple
# Default install location is ~/.local/share/deskflow-watchdog.
# If installed to a different prefix (e.g. /usr/local), update this path.
ExecStart=%h/.local/share/deskflow-watchdog/start-deskflow-watchdog.sh
Restart=on-failure
RestartSec=30
Environment="DISPLAY=:0"
Environment="WAYLAND_DISPLAY=wayland-0"

[Install]
WantedBy=default.target
